# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# -*- bash -*-
# shellcheck shell=bash

bashio::log.info "Starting NGINX..."
if [ -f /run/nginx ]; then
  rm -f /run/nginx
fi
mkdir -p /run/nginx

exec nginx -g 'daemon off;error_log /proc/1/fd/1 error;'