# services/nginx/run
#!/command/with-contenv bashio
# -*- bash -*-
# shellcheck shell=bash
mkdir /run/nginx
exec nginx -g 'daemon off;error_log /proc/1/fd/1 error;'